# roles/edge/vertexd/templates/vertexd.service
[Unit]
Description=Packetframe vertex rule manager

[Service]
User=root
ExecStart=/usr/local/bin/vertexd

[Install]
WantedBy=multi-user.target
